# config.in/loxprox-supervisor.ini
[program:loxprox]
command=loxprox -c /opt/nivos/loxprox/loxprox.yml
autostart=true
autorestart=true
stderr_logfile=/opt/nivos/loxprox/log/loxprox.err.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/opt/nivos/loxprox/log/loxprox.out.log
stdout_logfile_maxbytes=10MB

[inet_http_server]
port = 0.0.0.0:52080
username = loxprox
password = loxone
